# provisions/docker-builds/shared/config.supervisor/zookeeper.conf
[program:zookeeper]
user=zookeeper
command=/usr/lib/zookeeper/bin/zkServer.sh start-foreground
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/zookeeper/zookeeper.out
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
