# provisions/accumulo-1.6.5/supervisord/supervisord.conf
[supervisord]
nodaemon=false
pidfile=/tmp/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true


