# provisions/docker-builds/centos-base-ssh/supervisord-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
startretries=3

