# provisions/docker-builds/doop/bootstrap-node.conf
[program:bootstrap]
user=root
command=/usr/local/sbin/bootstrap-node.sh
priority=0
autostart=true
autorestart=false
startsecs=1
startretries=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/bootstrap
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
