# provisions/scidb/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:postgres]
command=/etc/init.d/postgresql restart